# deploy/hud-observer.service.example
[Unit]
Description=DoD HUD Observer backend (ingest + REST + Socket.IO)
After=network.target

[Service]
Type=simple
User=cadaver
WorkingDirectory=/opt/hud-observer
ExecStart=/usr/bin/npx ts-node --script-mode /opt/hud-observer/backend/src/app.ts
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=production
Environment=HUD_AUTH_KEY=REPLACE_WITH_SECRET
# Data server already has a KTP service on 8088 — move ingest to 9000.
Environment=HUD_INGEST_PORT=9000
StandardOutput=append:/var/log/hud-observer/backend.log
StandardError=append:/var/log/hud-observer/backend-err.log

[Install]
WantedBy=multi-user.target
